# infra/systemd/aer-verify-ilm.timer
# Phase 150 — daily ILM enforcement check. Runs after the 04:00 backup so it
# observes the post-prune state. The ILMCheckStale alert (>26h without a report)
# backstops a silently-missed run; ILMViolation fires if data lingers past TTL.
[Unit]
Description=Daily AĒR ILM TTL-enforcement check

[Timer]
OnCalendar=*-*-* 05:00:00
# Catch up a missed check if the box was down at 05:00.
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
